# cron/crontab
*/10 * * * * ~/current/jobs/enqueue-store-tweet-tasks.sh
