# jobs/crontab
*/10 * * * * ~/code/jobs/tweetstore/enqueue-store-tweet-tasks.sh
