# jobs/crontab
MAILTO="demoboxguard-laapbunny@yahoo.co.uk"
*/10 * * * * ~/code/jobs/tweetstore/enqueue-store-tweet-tasks.sh