# jobs/crontab
*/10 * * * * ~/jobs/tweetstore/enqueue-store-tweet-tasks.sh